# etc/docker/php/assets/supervisor/conf.d/nginx.conf
[program:nginx]
priority = 1
command = /usr/sbin/nginx
stdout_events_enabled=true
stderr_events_enabled=true