# etc/docker/php/assets/supervisor/conf.d/php-fpm.conf
[program:php-pm]
priority = 1
command = /app/vendor/php-pm/php-pm/bin/ppm start
autostart = true
user = www-data
redirect_stderr=true
stdout_logfile = /var/log/supervisor/php-pm.log
loglevel = warn