# etc/docker/php/assets/supervisor/conf.d/php-fpm.conf
[program:php-pm]
priority = 1
command = /usr/local/sbin/php-fpm -D
autostart = true
user = root
redirect_stderr=true
stdout_logfile = /var/log/supervisor/php-fpm.log
loglevel = warn